# crontab
#min   hour    day      month   weekday     command
*/5     *       *       *       *           root /usr/local/bin/drush -r /var/www/html cron >> /var/log/cron.log 2>&1