# crontab
#!/bin/sh

#min  hour    day     month   weekday   command
*/5     *       *       *       *         /usr/local/bin/drush -r /var/www/html --quiet cron